# references/config-templates/openclaw.service
[Unit]
Description=OpenClaw Gateway (AGENT_NAME_PLACEHOLDER)
Documentation=https://docs.openclaw.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw

# Environment file for API keys
EnvironmentFile=/etc/openclaw/env

# Working directory
WorkingDirectory=/home/openclaw/.openclaw

# Startup script
ExecStart=/usr/local/bin/openclaw-startup.sh

# Restart policy
Restart=on-failure
RestartSec=10

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/openclaw/.openclaw
ReadWritePaths=/tmp/openclaw
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw

[Install]
WantedBy=multi-user.target
